# prerendertron.service
[Unit]
Description=Rendertron cache populator
After=rendertron.service

[Service]
Type=oneshot
ExecStart=%h/rendertron/venv/bin/python %h/rendertron/prerendertron.py
